# mac/launchd/com.keepsake.migration.lens-more-than-human-time.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.keepsake.migration.lens-more-than-human-time</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/YOUR_USERNAME/keepsake-migration/mac/launchd/start-lens.sh</string>
        <string>more_than_human_time</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/YOUR_USERNAME/keepsake-migration/mac</string>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/Users/YOUR_USERNAME/keepsake-migration/mac/logs/lens-more-than-human-time.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/YOUR_USERNAME/keepsake-migration/mac/logs/lens-more-than-human-time.err</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin:/opt/homebrew/bin</string>
    </dict>

    <key>ThrottleInterval</key>
    <integer>30</integer>
</dict>
</plist>
