# mac/launchd/com.keepsake.custodian.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.keepsake.custodian</string>

  <key>UserName</key>
  <string>junyoo182</string>

  <key>ProgramArguments</key>
  <array>
    <string>/Users/junyoo182/keepsake-migration/mac/venv/bin/python</string>
    <string>-m</string>
    <string>helper.custodian.runner</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/junyoo182/keepsake-migration</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PYTHONUNBUFFERED</key>
    <string>1</string>
  </dict>

  <key>RunAtLoad</key>
  <true/>

  <key>KeepAlive</key>
  <dict>
    <key>SuccessfulExit</key>
    <false/>
  </dict>

  <key>ThrottleInterval</key>
  <integer>30</integer>

  <key>StandardOutPath</key>
  <string>/Users/junyoo182/keepsake-migration/helper/custodian/state/custodian.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/junyoo182/keepsake-migration/helper/custodian/state/custodian.error.log</string>
</dict>
</plist>
