# mac/launchd/com.keepsake.drift-br.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.keepsake.drift-br</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/Users/junyoo182/keepsake-migration/mac/launchd/start-drift-server.sh</string>
    <string>/Users/junyoo182/keepsake_drift_br</string>
    <string>8003</string>
  </array>

  <key>RunAtLoad</key>
  <true/>

  <key>KeepAlive</key>
  <dict>
    <key>SuccessfulExit</key>
    <false/>
  </dict>

  <key>ThrottleInterval</key>
  <integer>10</integer>

  <key>StandardOutPath</key>
  <string>/Users/junyoo182/keepsake-migration/mac/logs/drift-br.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/junyoo182/keepsake-migration/mac/logs/drift-br.error.log</string>
</dict>
</plist>
