# mac/launchd/com.keepsake.chrome-opener.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.keepsake.chrome-opener</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/YOUR_USERNAME/keepsake-migration/mac/scripts/open-admin-pages.sh</string>
    </array>

    <key>RunAtLoad</key>
    <true/>

    <!-- Run once only — do not restart after Chrome closes -->
    <key>KeepAlive</key>
    <false/>

    <key>StandardOutPath</key>
    <string>/Users/YOUR_USERNAME/keepsake-migration/mac/logs/chrome-opener.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/YOUR_USERNAME/keepsake-migration/mac/logs/chrome-opener.err</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin:/opt/homebrew/bin</string>
    </dict>
</dict>
</plist>
